# macchanger@wlan0.service

[Unit]
Description=Mac Spoof
Before=dhcpcd.service
 
[Service]
ExecStart=/bin/bash /usr/bin/randomize-mac.sh 
 
[Install]
WantedBy=network.target
